# services/rules_handler.service
[Unit]
Description=Rules handler for irrigation system
After=network.target
Requires=network.target

[Service]
EnvironmentFile=/etc/environment
ExecStartPre=/bin/sleep 40
ExecStart=/usr/bin/python3 /var/www/app/rules_handler/rules_handler.py

Restart=always
RestartSec=45
User=root

[Install]
WantedBy=multi-user.target
